# saab-control/saab.service
[Unit]
Description=SAAB Service
Restart=on-failure
RestartSec=2s

[Service]
ExecStart=python3 /usr/local/bin/SaabHeadUnit/saabCan.py

User=pi
Environment=DISPLAY=:0

