# saab-control/saab.service
[Unit]
Description=SAAB Service

[Service]
ExecStart=python3 /usr/local/bin/SaabHeadUnit/saab-control/saabCan.py

User=pi
Environment=DISPLAY=:0

